# pgbackuptimer/test-timer.service
[Unit]
Description=Test timer service

[Service]
Type=oneshot
# TODO update
ExecStart=/bin/bash %h/test-timer/test-timer.sh
